# com.valref.granola-upload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.valref.granola-upload</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>source ~/.zshrc &amp;&amp; cd ~/dev/valref &amp;&amp; echo "=== $(date) ===" &gt;&gt; logs/granola-upload.log 2&gt;&amp;1 &amp;&amp; yarn granola-upload &gt;&gt; logs/granola-upload.log 2&gt;&amp;1</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>TimeoutSeconds</key>
    <integer>300</integer>
</dict>
</plist>
